# garmin-sync/com.ncs.garminsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ncs.garminsync</string>

    <!-- OJO: apunta a ~/.local/share/ncs-garmin-sync/, NO a la carpeta del repo (Dropbox).
         launchd corre como daemon sin sesión de usuario y macOS (TCC) le niega el acceso
         a rutas dentro de ~/Library/CloudStorage/Dropbox/... con "Operation not permitted",
         aunque el mismo script funcione perfecto corrido a mano en una terminal.
         Tras editar garmin-sync/sync_garmin.py en el repo, hay que volver a copiarlo:
         cp garmin-sync/sync_garmin.py ~/.local/share/ncs-garmin-sync/sync_garmin.py -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/mauro/.local/share/ncs-garmin-sync/sync_garmin.py</string>
    </array>

    <!-- Corre todos los días a las 11:00. Si el Mac estaba dormido/apagado,
         launchd lo corre apenas despierta (no es garantizado, pero el script
         es idempotente así que un día salteado se recupera solo). -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>11</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/ncs-garmin-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ncs-garmin-sync.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
